# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/local/bin/service.sh
